# templates/default/tomcat_init.erb
#! /bin/bash

### BEGIN INIT INFO
# Provides:          tomcat
# Required-Start:    $network
# Required-Stop:     $network
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Apache Tomcat server
# Description:       Apache Tomcat server
### END INIT INFO

export JAVA_HOME="<%=@node[:tomcat6][:java_home]%>"

case $1 in
    start)
        <%=@node[:tomcat6][:tomcat_home]%>/bin/startup.sh
        exit 0
    ;;
    stop)
        <%=@node[:tomcat6][:tomcat_home]%>/bin/shutdown.sh
        exit 0
    ;;
    restart)
        $0 stop
        sleep 3
        $0 start
    ;;
    *)
    echo "Usage $0 {start|stop|restart}."
    exit 1
    ;;
esac
